# systemd/gpsd-watcher.service
[Unit]
Description=Persistent gpsd watcher - caches latest TPV for gpsloc
After=network.target gpsd.service
Requires=gpsd.service

[Service]
Type=simple
User=madvulcan
ExecStart=/usr/local/bin/gpsd-watcher
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
